# cron.d/rpgraph
*/5 * * * * www-data /opt/rpgraph/rpgraph_logger.sh
*/5 * * * * www-data /opt/rpgraph/rpgraph_buildchart.sh
